# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
; nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)


; Configured in /etc/proftpd.conf
[program:proftpd]
command=/sbin/proftpd --nodaemon

[program:s3fs-fuse]
command=bash -c '/usr/local/bin/s3fs -f -o public_bucket=1 -o enable_content_md5 -o ahbe_conf=/etc/ahbe_conf -o url=${S3_URL} ${BUCKET_NAME} ${STORAGE_PATH}'
